# installer/linux/sleepguard-monitor.service
[Unit]
Description=SleepGuard monitor (user session)
After=graphical-session.target

[Service]
Type=simple
ExecStart=@APP_BIN@ --background
Restart=always
RestartSec=5
Environment=SLEEPGUARD_CONFIG=@CONFIG_PATH@

[Install]
WantedBy=graphical-session.target
